# infrastructure/ansible/templates/backup-verify-quarterly.service.j2
[Unit]
Description=Backup verification (quarterly, integrity check)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=docker compose -f {{ project_root }}/infrastructure/docker/docker-compose.server.yml exec -T api node /app/apps/api/dist/scripts/verify-backups.js --mode=quarterly --max-age-days={{ backup_verify_quarterly_max_age_days | default(110) }} --max-download-mb={{ backup_verify_max_download_mb | default(50) }}
StandardOutput=journal
StandardError=journal

