# infrastructure/ansible/templates/kiosk-browser.service.j2
[Unit]
Description=Factory Kiosk Browser
After=graphical.target network-online.target
Wants=graphical.target

[Service]
Type=simple
User={{ ansible_user }}
Environment=DISPLAY=:0
Environment=GTK_USE_PORTAL=0
# キーリングを無効化（パスワードダイアログを表示しない）
Environment=GNOME_KEYRING_CONTROL=""
Environment=SSH_AUTH_SOCK=""
# NetworkManagerの認証ダイアログを抑制
Environment=NM_CLI_NO_TERSE=1
Environment=NM_CLI_NO_ASK_PASSWORD=1
# 日本語入力環境変数
Environment=GTK_IM_MODULE=ibus
Environment=XMODIFIERS=@im=ibus
Environment=QT_IM_MODULE=ibus
ExecStart=/usr/local/bin/kiosk-launch.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
