# infrastructure/ansible/templates/pi5-power-dispatcher.service.j2
[Unit]
Description=Pi5 power action dispatcher
After=network-online.target

[Service]
Type=oneshot
User=denkon5sd02
WorkingDirectory=/opt/RaspberryPiSystem_002/infrastructure/ansible
ExecStart=/usr/local/bin/pi5-power-dispatcher.sh /opt/RaspberryPiSystem_002/power-actions
StandardOutput=journal
StandardError=journal
