# infrastructure/ansible/templates/pi5-power-dispatcher.path.j2
[Unit]
Description=Watch power action queue

[Path]
PathExistsGlob=/opt/RaspberryPiSystem_002/power-actions/*.json
Unit=pi5-power-dispatcher.service

[Install]
WantedBy=multi-user.target
