# infrastructure/ansible/templates/private-pi5-stackchan-bridge.service.j2
[Unit]
Description=StackChan Private Bridge
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User={{ private_pi5_bridge_effective_user }}
WorkingDirectory={{ private_pi5_bridge_effective_dir }}
EnvironmentFile={{ private_pi5_bridge_env_path }}
ExecStart={{ private_pi5_bridge_venv_path }}/bin/python {{ private_pi5_bridge_effective_dir }}/bridge_server.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
