# infrastructure/ansible/templates/private-pi5-hermes-life-reminder.service.j2
[Unit]
Description=Hermes Life Pilot reminder notification dispatcher
After=network-online.target
Wants=network-online.target
ConditionPathExists={{ private_pi5_hermes_data_dir }}/.env
ConditionPathExists={{ private_pi5_hermes_data_dir }}/plugins/private-pi5-discord-task-bridge/life_reminder_scheduler.py

[Service]
Type=oneshot
User={{ private_pi5_hermes_user }}
Group={{ private_pi5_hermes_group }}
WorkingDirectory={{ private_pi5_hermes_data_dir }}/plugins/private-pi5-discord-task-bridge
Environment=HOME={{ private_pi5_hermes_home }}
Environment=LIFE_PILOT_STORAGE_ROOT={{ private_pi5_hermes_life_data_dir }}
EnvironmentFile={{ private_pi5_hermes_data_dir }}/.env
ExecStart=/usr/bin/python3 {{ private_pi5_hermes_data_dir }}/plugins/private-pi5-discord-task-bridge/life_reminder_scheduler.py --storage-root {{ private_pi5_hermes_life_data_dir }}
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
