# infrastructure/ansible/templates/security-monitor.service.j2
[Unit]
Description=Security log monitor
ConditionPathExists={{ alert_script_path }}

[Service]
Type=oneshot
ExecStart=/usr/local/bin/security-monitor.sh

