# infrastructure/ansible/templates/signage-daily-reboot.timer.j2
[Unit]
Description=Daily reboot timer for Pi3 signage stability
After=network-online.target

[Timer]
# 毎日3時（深夜）に再起動
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target

