# infrastructure/ansible/roles/signage/templates/signage-lite-update.timer.j2
[Unit]
Description=Update signage image periodically
After=network-online.target

[Timer]
OnBootSec=1min
OnUnitActiveSec={{ signage_update_interval | default(30) }}s
# systemd timerのデフォルトAccuracySec(1min)だと、30秒タイマーでも1分単位に丸められ
# 「更新が30秒おきにならない → ページが飛ぶ」事象につながるため、精度を明示する。
AccuracySec=1s
RandomizedDelaySec=0

[Install]
WantedBy=timers.target

