# infrastructure/ansible/roles/signage/templates/signage-lite.service.j2
[Unit]
Description=Digital Signage Lite (feh-based)
After=graphical.target
Wants=graphical.target
# ネットワーク接続を待たない（オフライン時でも既存画像を表示）
# network-online.target への依存を削除

[Service]
Type=simple
User={{ ansible_user }}
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/{{ ansible_user }}/.Xauthority

# DISPLAYが利用可能になるまで待機（最大30秒）
ExecStartPre=/bin/bash -c 'for i in {1..30}; do [ -n "$DISPLAY" ] && xset q >/dev/null 2>&1 && break || sleep 1; done'

# 画像表示（内部で画像取得を試行）
ExecStart=/usr/local/bin/signage-display.sh

# ネットワーク遮断時でもサービスを再起動し続ける
# 既存画像があれば表示し続ける
Restart=always
RestartSec=10

# 暴走防止: 5分間に5回以上失敗した場合は10分間停止
StartLimitIntervalSec=300
StartLimitBurst=5
StartLimitAction=none

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target

