# infrastructure/ansible/roles/client/templates/torque-bluetooth-adapter@.service.j2
[Unit]
Description=Fail-closed initialization of exact external torque Bluetooth adapter %I
After=bluetooth.service systemd-rfkill.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/torque-bluetooth-adapter --power-off %I
TimeoutStartSec=90
TimeoutStopSec=10
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
